# packages/autostart/org.marishoku.system-ready.desktop
[Desktop Entry]
Type=Application
Name=MARISHOKU System Ready
Comment=Show the URA login status dialog
Exec=marishoku-system-ready
Icon=marishoku-heart
OnlyShowIn=KDE;
Terminal=false
X-KDE-autostart-phase=2
X-KDE-StartupNotify=false
